# debian/spark-modem-watchdog.service
[Unit]
Description=Spark modem watchdog (Sierra EM7421 + Soliton Zao)
Documentation=https://github.com/spark/spark-modem-watchdog
After=network-online.target zao-infra-ctrl.service
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main

# FR-60 + B-03 belt-and-suspenders B: smoke-test on every start.
ExecStartPre=/opt/spark-modem-watchdog/libexec/postinst_smoke_test.sh

# Phase 1 placeholder: the real daemon entry point lands in Phase 2.
# For now, ExecStart= is a no-op that immediately reports READY and exits 0
# so `systemctl start` succeeds end-to-end after a clean install.
# Phase 2 will replace this with: ExecStart=/opt/spark-modem-watchdog/bin/spark-modem-watchdog
ExecStart=/opt/spark-modem-watchdog/python/bin/python3.12 -c \
  'import sdnotify; sdnotify.SystemdNotifier().notify("READY=1"); import time; time.sleep(0.1)'

Restart=on-failure
RestartSec=5s

# Phase 1: hardened baseline. Phase 2 relaxes CapabilityBoundingSet to add
# CAP_NET_ADMIN for pyroute2; Phase 4 adds CAP_SYS_ADMIN for usb_reset/driver_reset.
User=spark-modem-watchdog
Group=spark-modem-watchdog
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
# Daemon needs: AF_UNIX (Prometheus UDS), AF_INET/AF_INET6 (webhooks),
# AF_NETLINK (pyroute2 rtnetlink).
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Phase 1 needs no capabilities (placeholder ExecStart only notifies READY=1).
# Phase 2 will add CAP_NET_ADMIN; Phase 4 will add CAP_SYS_ADMIN.
CapabilityBoundingSet=
AmbientCapabilities=

# systemd-managed directories: auto-created with correct ownership, auto-cleaned.
RuntimeDirectory=spark-modem-watchdog
RuntimeDirectoryMode=0750
StateDirectory=spark-modem-watchdog
StateDirectoryMode=0750
LogsDirectory=spark-modem-watchdog
LogsDirectoryMode=0750
ConfigurationDirectory=spark-modem-watchdog

ReadWritePaths=/var/lib/spark-modem-watchdog /var/log/spark-modem-watchdog /run/spark-modem-watchdog

# ADR-0011: HMAC secret loaded via LoadCredential= (never baked into the unit or env files).
# The credential file path is /etc/credstore/spark-modem-watchdog.hmac-secret
# or set via systemd's LoadCredential= mechanism at deploy time.
# Phase 2 implementation reads via CREDENTIALS_DIRECTORY env var set by systemd.
LoadCredential=spark-modem-watchdog.hmac-secret:/etc/spark-modem-watchdog/hmac-secret

[Install]
WantedBy=multi-user.target
