# debian/spark-modem-watchdog.service
[Unit]
Description=Spark modem watchdog (Sierra EM7421 + Soliton Zao)
Documentation=https://github.com/spark/spark-modem-watchdog
After=network-online.target zao-infra-ctrl.service
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main

# FR-60 + B-03 belt-and-suspenders B: smoke-test on every start.
ExecStartPre=/opt/spark-modem-watchdog/libexec/postinst_smoke_test.sh

# U-05: pre-flight Settings validate before main process. Catches bad
# config rollouts BEFORE the main daemon process boots, so a bad
# config doesn't trip StartLimitBurst (PITFALLS §4.2).
ExecStartPre=/opt/spark-modem-watchdog/bin/spark-modem ctl config-check

# Phase 3 wires the real entry point. The wrapper script lives in
# /opt/spark-modem-watchdog/bin/ and is provided by the .deb postinst
# (Phase 4 follow-up: verify wrapper exists; current .deb scaffold may
# need an additional ExecStart shim).
ExecStart=/opt/spark-modem-watchdog/bin/spark-modem-watchdog

# U-02: clean SIGTERM exit must not trigger restart. Default StartLimit
# (5 restarts in 10s) bricks fleet on bad config rollout (PITFALLS §4.2).
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=20
TimeoutStopSec=10s
KillMode=mixed

# U-04: WatchdogSec=90s. Daemon kicks WATCHDOG=1 at cycle-end (PITFALLS
# §4.1); a stuck mid-cycle triggers systemd-managed restart at the 90s mark.
# Phase 4 HIL stress tests verify the watchdog kicks even during cycles
# approaching the 10s P99 budget (NFR-1).
WatchdogSec=90s

# NFR-30: daemon runs as root; Phase 3 needs CAP_NET_ADMIN on udev/pyroute2,
# Phase 4 will exercise CAP_SYS_ADMIN/CAP_SYS_MODULE for usb_reset/driver_reset.
# NoNewPrivileges=yes pins the safety floor; sandboxing below provides
# defence-in-depth without breaking LoadCredential= on systemd 245.
User=root
Group=root
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
# U-03: allow netns + mnt for `ip netns exec` and the prom UDS bind.
RestrictNamespaces=net mnt
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
# Daemon needs: AF_UNIX (Prometheus UDS), AF_INET/AF_INET6 (webhooks),
# AF_NETLINK (pyroute2 rtnetlink + pyudev kobject_uevent).
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# U-03 sandboxing notes:
#   - No PrivateMounts= (PITFALLS §4.3 — incompatible with LoadCredential=
#     on systemd 245 / Ubuntu 20.04).
#   - No PrivateTmp= (drops for LoadCredential compat + /run visibility for
#     a sibling `spark-modem ctl` mutator that needs the same flock files).
#   - No PrivateDevices= (kmsg producer needs read access to /dev/kmsg).

# U-01: CapabilityBoundingSet — Phase 4-forward. Phase 3 doesn't USE
# CAP_SYS_MODULE (driver reset lands in Phase 4) but ships it preallocated
# so the unit-file is stable across the Phase 3 → Phase 4 boundary.
# Single unit-file edit, no mid-rollout edits in Phase 4.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_ADMIN CAP_SYS_MODULE CAP_DAC_READ_SEARCH
AmbientCapabilities=

# systemd-managed directories: auto-created with correct ownership, auto-cleaned.
RuntimeDirectory=spark-modem-watchdog
RuntimeDirectoryMode=0750
# U-03 / PITFALLS §4.4: RuntimeDirectoryPreserve=yes is LOAD-BEARING.
# Without it, systemd cleans /run/spark-modem-watchdog/ on stop, taking out
# the PID lock + clean-shutdown marker + state.lock + modem-*.lock + metrics.sock.
RuntimeDirectoryPreserve=yes
StateDirectory=spark-modem-watchdog
StateDirectoryMode=0750
LogsDirectory=spark-modem-watchdog
LogsDirectoryMode=0750
ConfigurationDirectory=spark-modem-watchdog

ReadWritePaths=/var/lib/spark-modem-watchdog /var/log/spark-modem-watchdog /run/spark-modem-watchdog

# ADR-0011: HMAC secret loaded via LoadCredential= (never baked into the unit or env files).
# The credential file path is /etc/credstore/spark-modem-watchdog.hmac-secret
# or set via systemd's LoadCredential= mechanism at deploy time.
# Phase 2 implementation reads via CREDENTIALS_DIRECTORY env var set by systemd.
LoadCredential=spark-modem-watchdog.hmac-secret:/etc/spark-modem-watchdog/hmac-secret

[Install]
WantedBy=multi-user.target
